# src/modules/adsb-feeder/filesystem/root/usr/lib/systemd/system/adsb-push-prometheus-metrics.timer
[Unit]
Description=Push Prometheus metrics every 3 minutes
Requisite=adsb-setup.service
After=adsb-setup.service

[Timer]
OnCalendar=*-*-* *:0/3:0
AccuracySec=5s

[Install]
WantedBy=multi-user.target
