# src/modules/adsb-feeder/filesystem/root/usr/lib/systemd/system/adsb-expand-rootfs.service
[Unit]
Description=Expand the root file system

[Service]
Type=oneshot
ExecStart=/opt/adsb/expand-rootfs.bash

[Install]
WantedBy=sysinit.target
